# deploy/polymomentum-healthcheck.service
[Unit]
Description=PolyMomentum healthcheck (one-shot, fired by .timer)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/polymomentum/env
ExecStart=/opt/polymomentum/healthcheck.sh
StandardOutput=journal
StandardError=journal
