# deploy/polycrossarb-rust.service
[Unit]
Description=PolyCrossArb rust pipeline
After=network-online.target

[Service]
Type=simple
User=polycrossarb
WorkingDirectory=/opt/polycrossarb
EnvironmentFile=/opt/polycrossarb/.env
ExecStart=/opt/polycrossarb/polycrossarb-engine
Restart=on-failure
Environment=CLOB_DIRECT=1
LimitNOFILE=65536
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
