# deploy/polycrossarb-candle.service
[Unit]
Description=PolyCrossArb candle pipeline
After=network-online.target
Wants=network-online.target

[Service]
Environment=PYTHONUNBUFFERED=1
Environment=UV_CACHE_DIR=/tmp/polycrossarb/uv-cache
Type=simple
User=polycrossarb
WorkingDirectory=/opt/polycrossarb/current
EnvironmentFile=/etc/polycrossarb/env
ExecStart=/usr/local/bin/uv run python scripts/run_production.py --mode paper
ExecStop=/bin/kill -SIGTERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/polycrossarb /tmp/polycrossarb

[Install]
WantedBy=multi-user.target
