# deploy/polycrossarb-healthcheck.service
[Unit]
Description=PolyCrossArb healthcheck (one-shot, fired by .timer)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/polycrossarb/env
ExecStart=/opt/polycrossarb/current/deploy/healthcheck.sh
StandardOutput=journal
StandardError=journal
